# templates/vault.upstart.erb
<% require 'shellwords' -%>
# Vault Server (Upstart unit)
description "Vault Server"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

respawn
respawn limit 10 10
kill timeout 10
<% if !scope.lookupvar('vault::oom_score').nil? or !scope.lookupvar('vault::oom_score').to_s == 'undef' -%>
oom score <%= scope.lookupvar('vault::oom_score') %>
<% end -%>

env GOMAXPROCS=<%= scope.lookupvar("::processorcount") %>
env CONFIG='<%= scope.lookupvar('vault::config_dir') %>/config.json'

script
    # read settings like GOMAXPROCS from "/etc/default/vault", if available.
    [ -e "/etc/default/vault" ] && . "/etc/default/vault"
    exec vault server -config "$CONFIG" <%= scope.lookupvar('vault::extra_options').shelljoin %>
end script
